# misc/upstart.conf
description "CertProxy Daemon"

start on (started networking)
stop on runlevel [!2345]

respawn
respawn limit 5 5

console none

exec /usr/local/bin/certproxy --config /etc/moneta/certproxy.yml
